# deploy/systemd/campus-hub.service
[Unit]
Description=Campus Hub Next.js App
After=network.target

[Service]
Type=simple
Environment=NODE_ENV=production
Environment=PORT=3000
WorkingDirectory=/var/www/Campus_Hub
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
User=campus
Group=campus

[Install]
WantedBy=multi-user.target