# deploy/systemd/campus-hub-api.service
[Unit]
Description=Campus Hub Local API (server.js)
After=network.target

[Service]
Type=simple
Environment=NODE_ENV=production
WorkingDirectory=/var/www/Campus_Hub
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target


